# infra/geodns/systemd/majlesyar-geoipupdate.service
[Unit]
Description=Refresh MaxMind GeoLite2 Country database
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/geoipupdate
